# systemd/uftp-fpv-client.container
# Place in: /etc/containers/systemd/   (system)
#        or: ~/.config/containers/systemd/  (user, rootless)
# Then run: systemctl daemon-reload
#           systemctl start uftp-fpv-client

[Unit]
Description=FPV Camera Client — captures frames and sends via UFTP
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/aeronext/uftp-fpv-client:latest

# Camera device passthrough
AddDevice=/dev/video0

# Network (host mode ensures UDP can reach the receiver directly)
Network=host

# 環境変数はファイルから読み込む
# cp env/client.env /etc/uftp-fpv/client.env で配置後、内容を編集してください
EnvironmentFile=/etc/uftp-fpv/client.env

# Restart on failure, not on success (i.e. SIGTERM on stop is intentional)
[Service]
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=default.target
