# systemd/uftp-fpv-receiver.container
# Place in: /etc/containers/systemd/   (system)
#        or: ~/.config/containers/systemd/  (user, rootless)
# Then run: systemctl daemon-reload
#           systemctl start uftp-fpv-receiver

[Unit]
Description=UFTP Receiver — accepts incoming camera frames from drone/client
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/aeronext/uftp-fpv-receiver:latest

# Host networking so UFTP multicast/unicast UDP 1044 is reachable
Network=host

# Shared named volume with the web viewer
Volume=fpv-images.volume:/data/images:z

# 環境変数はファイルから読み込む
# cp env/receiver.env /etc/uftp-fpv/receiver.env で配置後、内容を編集してください
EnvironmentFile=/etc/uftp-fpv/receiver.env

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
